# config/etc/supervisord.conf
[program:scheduler]
command=python /app_django/scheduler.py
numprocs=1
autostart=true
autorestart=true